# sys-power/fanctl/files/fanctl.initd
#!/sbin/openrc-run

command=/usr/bin/fanctl
command_args="-c /etc/fanctl.yaml &"
